# ansible/roles/shard_management/templates/oci-shard-sync.timer.j2
[Unit]
Description=Run OCI Shard Sync every 10 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target

